# systemd/youtube-live-translate.service
[Unit]
Description=Local YouTube Live Translate Service
After=network.target

[Service]
Type=simple
User=<your-user>
WorkingDirectory=/absolute/path/to/youtube-live-translate
EnvironmentFile=-/absolute/path/to/youtube-live-translate/.env
ExecStart=/absolute/path/to/node /absolute/path/to/youtube-live-translate/dist/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
